# opennet/packages/on-captive-portal/files/etc/hotplug.d/iface/80-on-captive-portal-discovery
../../../usr/lib/opennet/events/on-captive-portal-discovery.sh