# opennet/packages/on-captive-portal/files/etc/hotplug.d/iface/on-captive-portal-hotplug.sh
../../../usr/lib/opennet/events/on-captive-portal-hotplug.sh