# opennet/packages/on-nodogsplash/files/etc/hotplug.d/iface/80-on-nodogsplash-discovery
../../../usr/lib/opennet/events/on-nodogsplash-discovery.sh