# opennet/packages/on-olsr2/files/etc/hotplug.d/iface/45-on-olsr2
#!/bin/sh

case "$ACTION" in
	ifup)
		# leicht verzoegert durchfuehren, um Hotplug-Schleifen zu vermeiden
		echo "on-function update_olsr2_interfaces; on-function init_policy_routing_ipv6" | on-function schedule_task
		;;
esac
