# opennet/packages/on-openvpn/files/etc/init.d/on_openvpn_config
#!/bin/sh /etc/rc.common
#
# Opennet Firmware
# 
# Copyright 2010 Rene Ejury <opennet@absorb.it>
# 
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
# 
#   http://www.apache.org/licenses/LICENSE-2.0
# 

START=21

start() {
	. "${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh"

	# create uci node
	[ -e /etc/config/on-openvpn ] || touch /etc/config/on-openvpn

	# firewall settings
	if [ -z "$(uci_get network.on_vpn)" ]; then
		local batch
		
		# add new network to configuration (to be recognized by olsrd)
		append batch "set network.on_vpn=interface${N}"
		append batch "set network.on_vpn.proto=none${N}"
		append batch "set network.on_vpn.ifname=tun0${N}"

		append batch "set firewall.zone_on_vpn=zone${N}"
		append batch "set firewall.zone_on_vpn.name=on_vpn${N}"
		append batch "set firewall.zone_on_vpn.network=on_vpn${N}"
		append batch "set firewall.zone_on_vpn.forward=ACCEPT${N}"
		append batch "set firewall.zone_on_vpn.input=ACCEPT${N}"
		append batch "set firewall.zone_on_vpn.output=ACCEPT${N}"
		append batch "set firewall.zone_on_vpn.masq=1${N}"

		echo "$batch${N}commit network${N}" | uci -q batch

		section=$(uci add firewall forwarding)
		uci set firewall.$section.src='local'
		uci set firewall.$section.dest='on_vpn'

		uci commit firewall
	fi

	# openvpn connection
	prefix=openvpn.opennet_user
	if [ -z "$(uci_get "$prefix")" ]; then
		uci set "${prefix}=openvpn"
		uci set "${prefix}.client=1"
		uci set "${prefix}.enable=1"
		uci set "${prefix}.proto=udp"
		uci set "${prefix}.dev=tun"
		uci set "${prefix}.persist_key=1"
		uci set "${prefix}.persist_tun=1"
		uci set "${prefix}.ca=/etc/openvpn/opennet_user/opennet-ca.crt"
		uci set "${prefix}.cert=/etc/openvpn/opennet_user/on_aps.crt"
		uci set "${prefix}.key=/etc/openvpn/opennet_user/on_aps.key"
		uci set "${prefix}.nobind=1"
		uci set "${prefix}.comp_lzo=yes"
		uci set "${prefix}.up=/etc/openvpn/opennet_user/opennet_up.sh"
		uci set "${prefix}.down=/etc/openvpn/opennet_user/opennet_down.sh"
		uci set "${prefix}.ns_cert_type=server"
		uci set "${prefix}.rport=1600"
		uci set "${prefix}.script_security=2 execve"
		uci add_list "${prefix}.remote=remote"
		uci set "${prefix}.remote=192.168.0.254"
		uci commit "$prefix"
	fi

	# general on-openvpn settings (gateways, preferred order, ...)
	prefix=on-openvpn.gateways
	if [ -z "$(uci_get "$prefix")" ]; then
		uci set "${prefix}=gateways"
		uci set "${prefix}.better_gw=0"
		uci set "${prefix}.current_gw="
		uci set "${prefix}.searchmask"='$1 ~ "^192\\.168\\.0\\.[0-9]+$" && $1 != "192.168.0.0"'
		uci set "${prefix}.vpn_recheck_age=20"
		uci set "${prefix}.vpn_nonworking_timeout=5"
		uci set "${prefix}.vpn_bettergateway_timeout=5"
		uci set "${prefix}.autosearch=on"
		uci set "${prefix}.vpn_sort_criteria=etx"
		uci commit "$prefix"
	fi

	# openssl certificate settings
	prefix=on-openvpn.openssl
	if [ -z "$(uci_get "$prefix")" ]; then
		uci set "${prefix}=openssl"
		uci set "${prefix}.countryName=de"
		uci set "${prefix}.provinceName=Mecklenburg-Vorpommern"
		uci set "${prefix}.localityName=Rostock"
		uci set "${prefix}.organizationalUnitName=users"
		uci set "${prefix}.days=3650"
		uci commit "$prefix"
	fi

    exit 0
}
