# opennet/packages/on-core/files/etc/cron.hourly/olsr_remove_old_services
#!/bin/sh

on-function remove_old_olsr_services
