# opennet/packages/on-usergw/files/etc/cron.daily/on_mesh_gateway_speed_test
#!/bin/sh

. "${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh"

is_on_module_installed_and_enabled "on-usergw" || exit 0

update_wan_gateway_speeds() {
	get_services "${RELAYABLE_SERVICE_PREFIX}-gw" \
			| filter_reachable_services \
			| filter_enabled_services \
			| while read service_name; do
		update_public_gateway_speed_estimation "$service_name"
	done
}

# alle APs fuehren den daily-cronjob zeitgleich aus - wir moechten aber nicht gleichzeitig die UGW-Server belasten
run_delayed_in_background "$(get_random 7200)" update_wan_gateway_speeds
