# opennet/packages/on-certificates/files/etc/cron.hourly/on_certificates_update_after_boot
#!/bin/sh

set -eu

# ermittle die System-Laufzeit
uptime_minutes=$(cut -f 1 -d . /proc/uptime | awk '{ print int($1 / 60) }')

# Innerhalb der ersten zwei Stunden wollen wir das Skript immer ausfuehren.
# Dies ist besonders sinnvoll fuer APs, die nur gelegentlich und kurz angeschaltet werden.
[ "$uptime_minutes" -gt 120 ] && exit 0

/usr/sbin/on_certificates_update
