# opennet/packages/on-openvpn/files/etc/cron.hourly/on_update_traceroute
#!/bin/sh

echo "on-function update_traceroute_gw_cache" | on-function schedule_task
