# opennet/packages/on-core/files/etc/init.d/on-core
#!/bin/sh /etc/rc.common
#
# Ausführung von opennet-relevanten Boot-Aktionen.
#

USE_PROCD=1
START=99
STOP=2


boot() {
	# sorge dafür, dass alle Skripte im opennet-boot-Verzeichnis zeitnah ausgeführt werden
	on-function schedule_parts /usr/lib/opennet/boot
}
