# opennet/packages/on-core/files/etc/init.d/on-core
#!/bin/sh /etc/rc.common
#
# Ausf端hrung von opennet-relevanten Boot-Aktionen.
#

USE_PROCD=1
START=99
STOP=2


boot() {
	# wir d端rfen die Einbindung nicht im globalen Namensraum vornehmen, da rc.common nicht im strikten Modus funktioniert
	. "${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh"
	local fname
	# sorge daf端r, dass alle Skripte im opennet-boot-Verzeichnis zeitnah ausgef端hrt werden
	find /usr/lib/opennet/boot -type f -perm +100 | while read fname; do
		echo "$fname" | schedule_task
	done
}
