# opennet/packages/on-usergw/files/etc/init.d/on_usergw
#!/bin/sh /etc/rc.common

START=31

start() {
	# stelle sicher, dass immer mindestens die beiden Standard-UGWs vorhanden sind
	on-function update_ugw_services
}

