# opennet/packages/on-usergw/files/etc/cron.hourly/on_service_relay
#!/bin/sh

set -eu

on-function update_service_relay_status
