# opennet/packages/on-tinc/files/etc/init.d/on_tinc
#!/bin/sh /etc/rc.common

. "${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh"

START=42
on_id=$(uci_get on-core.settings.on_id)

NETWORK_NAME="oni-media"
# you have to configure your tinc network address here.
# Please ask the other tinc-members to prevent collisions
# TINC_ADDRESS="10.42.1.X"
TINC_ADDRESS=""
TINC_NETMASK="255.255.255.0"
TINC_PORT=655
# you might like to change your node name here
# The name should consist only of alfanumeric and underscore characters (a-z, A-Z, 0-9 and _)
NODE_NAME="oni_${on_id/./_}"
FIREWALL_FILE=/etc/firewall.tinc


### create a new tinc config, generate keys and setup initial configuration
create_tinc_config() {
	local index
	local section
	local tinc_dir=/etc/tinc/$NETWORK_NAME
	mkdir -p "$tinc_dir/hosts"
	cat >"$tinc_dir/tinc.conf" << EOF
Mode=switch
Name=$NODE_NAME
Port=$TINC_PORT
# change this to connect and add key to hosts dir
#ConnectTo = <name>
EOF
	cat > "$tinc_dir/hosts/$NODE_NAME" << EOF
Address=$TINC_ADDRESS $TINC_PORT
Compression=9
EOF
	tincd -n "$NETWORK_NAME" --generate-keys=4096 </dev/null

	cat > "$tinc_dir/tinc-up" << EOF
#!/bin/sh
ifconfig "\$INTERFACE" $TINC_ADDRESS netmask $TINC_NETMASK"
EOF
	chmod +x "$tinc_dir/tinc-up"

	# Netzwerke hinzufuegen
	cat >>/etc/firewall.tinc << EOF
iptables -I input_rule -p udp --dport $TINC_PORT -j ACCEPT
iptables -I input_rule -p tcp --dport $TINC_PORT -j ACCEPT
iptables -I input_rule -i $NETWORK_NAME -j ACCEPT
iptables -I output_rule -o $NETWORK_NAME -j ACCEPT
EOF

	index=0
	while [ -n "$(uci_get firewall.@include[$index])" ]; do
		[ "$(uci_get firewall.@include[$index].path)" = "$FIREWALL_FILE" ] && return
		: $((index++))
	done
	# die Datei wurde noch nicht konfiguriert
	section=$(uci add firewall include)
	uci set "firewall.${section}.path=$FIREWALL_FILE"
	uci commit firewall
}

start() {
	local on_id_preset=$(get_on_core_default on_id_preset)
	[ -z "$on_id" -o "$on_id" = "$on_id_preset" ] && msg_info "not starting, on_id not configured" && return
	[ -z "$TINC_ADDRESS" ] && msg_info "not starting, TINC_ADDRESS in /etc/init.d/on_tinc not configured" && return
	[ ! -d "/etc/tinc/$NETWORK_NAME" ] && create_tinc_config
	tincd -n "$NETWORK_NAME"
}

stop() {
	tincd -n "$NETWORK_NAME" -k
}

restart() {
	stop
	start
}

