# opennet/packages/on-core/files/etc/uci-defaults/20-on-migrations
../../usr/lib/opennet/boot/on-migrations