# opennet/packages/on-core/files/etc/uci-defaults/45-on-version-banner
../../usr/lib/opennet/init/on-version-banner