# opennet/packages/on-core/files/etc/uci-defaults/60-on-hostapd-channel-whitelist
../../usr/lib/opennet/init/on-hostapd-channel-whitelist