# opennet/packages/on-core/files/etc/uci-defaults/70-on-openwrt-fixes
../../usr/lib/opennet/init/on-openwrt-fixes