# opennet/packages/on-core/files/etc/uci-defaults/on-discard-download-url
../../usr/lib/opennet/init/on-discard-download-url