# opennet/packages/on-core/files/etc/uci-defaults/on-migrations
../../usr/lib/opennet/init/on-migrations