# opennet/packages/on-core/files/etc/uci-defaults/on-migrations
#!/bin/sh
#
# Dieses Skript enthaelt Funktionen, die fuer Updates altes Firmware-Versionen notwendig sind.
# Alle Funktionen muessen idempotent sein.
#
# Dieses Skript wird nur ein einziges Mal nach einem Upgrade oder der Erstinstallation ausgefuehrt:
#   http://wiki.openwrt.org/doc/uci#defaults
#


# lieber strikte Fehlerpruefungen - wir wollen nix kaputtmachen, wenn es nicht den Erwartungen entspricht
set -eu

. "${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh"


# bis Version v0.4-5: openvpn.opennet_user.comp_lzo=1
# seit Version v0.5 muss die Einstellung einen der folgenden Werte haben: yes/no/adaptive
# Wir ersetzen 1 durch "yes" und 0 durch "no"
# Status in v0.4-5:
#   ~# uci show|grep lzo
#   on-usergw.opennet_ugw.comp_lzo=1
#   openvpn.opennet_user.comp_lzo=1
#   openvpn.opennet_ugw_erina_on_i_de.comp_lzo=1
#   openvpn.opennet_ugw_subaru_on_i_de.comp_lzo=1
coerce_openvpn_comp_lzo() {
	local new_value
	local key
	(uci show openvpn; uci show on-usergw) | grep "\.comp_lzo=[01]$" | while read line; do
		new_value=no
		echo "$line" | grep -q "1$" && new_value=yes
		key=$(echo "$line" | cut -f 1 -d =)
		uci set "$key=$new_value"
	done
	if [ -n "$(uci changes openvpn)" -o -n "$(uci changes on-usergw)" ]; then
		msg_info "MIGRATION: coerce_openvpn_comp_lzo"
		uci commit on-usergw
		uci commit openvpn
		/etc/init.d/openvpn restart
	fi
}


# bis Version v0.4-5: die Opennet-Firmware hat die /etc/passwd durch einen Symlink auf /etc/etc_preset/passwd ersetzt
# Bei einem Update wird das symlink-Ziel ersetzt und somit gibt es keine Nutzerdatenbank mehr.
# Dies verhindert jeden telnet/ssh-Login-Versuch. Lediglich das Web-Interface ist nutzbar.
fix_passwd_broken_symlink() {
	local target=/etc/passwd
	if [ -h "$target" -a ! -e "$target" ]; then
		msg_info "MIGRATION: fix_passwd_broken_symlink"
		rm "$target"
		# ein huebscheres here-Document mit Tabulator-Bereinigung ("<<-") funktioniert leider nicht mit busybox
		cat >"$target" << EOF
root:x:0:0:root:/root:/bin/ash
daemon:*:1:1:daemon:/var:/bin/false
ftp:*:55:55:ftp:/home/ftp:/bin/false
network:*:101:101:network:/var:/bin/false
nobody:*:65534:65534:nobody:/var:/bin/false
EOF
		# Ein paar Dienste schlugen aufgrund der fehlenden Nutzerdatenbank fehl.
		# Ein reboot waere schoen - aber kann zukuenftig eventuell irgendwann zu einer Schleife fuehren.
		# Also: manuell einzelne Dienste neu starten.
		/etc/init.d/dnsmasq restart
	fi
}


# bis Version v0.4-5: "firewall reload" fuehrte auch die "include"-Dateien aus (z.B. /etc/firewall.opennet)
# Ab Version v0.5 passiert dies nur noch, falls der include-Eintrag via uci mit dem Parameter 'reload' markiert wird.
# Beispiel (v0.4-5):
#  uci show firewall | grep @include
#  firewall.@include[0]=include
#  firewall.@include[0].path=/etc/firewall.opennet
#  firewall.@include[1]=include
#  firewall.@include[1].path=/etc/firewall.user
fix_firewall_reload() {
	local key
	local value
	local reload_key
	local old_value
	uci show firewall | grep "@include\[[0-9]\+\]\.path=" | while read line; do
		key=$(echo "$line" | cut -f 1 -d =)
		value=$(echo "$line" | cut -f 2- -d =)
		[ "$value" != "/etc/firewall.opennet" ] && continue
		reload_key=${key%.path}.reload
		old_value=$(uci_get "$reload_key")
		[ -z "$old_value" ] && old_value=0
		if uci_is_false "$old_value"; then
			msg_info "MIGRATION: fix_firewall_reload"
			uci set "$reload_key=1"
			uci commit firewall
			# Firewall neustarten; keine unnoetigen Ausgaben (siehe /etc/init.d/firewall)
			QUIET=-q /etc/init.d/firewall reload
			break
		fi
	done
}


coerce_openvpn_comp_lzo
fix_passwd_broken_symlink
fix_firewall_reload

