# opennet/packages/on-core/files/etc/uci-defaults/on-olsr-setup
#!/bin/sh
#
# Setzen aller fuer Opennet relevanten olsr-Einstellungen.
#
# Dieses Skript wird nur ein einziges Mal nach einem Upgrade oder der Erstinstallation ausgefuehrt:
#   http://wiki.openwrt.org/doc/uci#defaults
#

set -eu

. "${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh"


# erwuenschte Einstellungen setzen
enable_ondataservice
enable_nameservice
disable_missing_olsr_modules
olsr_set_main_ip
olsr_set_routing_tables olsrd olsrd-default

uci commit olsrd
/etc/init.d/olsrd restart >/dev/null || true

