# opennet/packages/on-core/files/etc/uci-defaults/on-olsr-setup
../../usr/lib/opennet/init/on-olsr-setup