# opennet/packages/on-core/files/etc/uci-defaults/on-password
#!/bin/sh
#
# Dieses Skript wird nur ein einziges Mal nach einem Upgrade oder der Erstinstallation ausgefuehrt:
#   http://wiki.openwrt.org/doc/uci#defaults
#

set -eu

# nach jedem Update wird das Passwort neu gesetzt
echo "root:admin" | chpasswd

