# opennet/packages/on-openvpn-v6/files/etc/uci-defaults/60-on-openvpn-v6
../../usr/lib/opennet/init/on-openvpn-v6