# opennet/packages/on-core/files/etc/hotplug.d/iface/40-on-core
#!/bin/sh

case "$ACTION" in
        ifup)
		on-function update_olsr_interfaces
		on-function update_opennet_zone_masquerading
		on-function initialize_olsrd_policy_routing
		;;
esac

