# opennet/packages/on-core/files/etc/hotplug.d/iface/40-on-core
#!/bin/sh

case "$ACTION" in
        ifup)
		# alle Aktionen leicht verzoegert durchfuehren, um Hotplug-Schleifen zu vermeiden
		echo ". ${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh; update_olsr_interfaces; update_opennet_zone_masquerading; initialize_olsrd_policy_routing" | on-function schedule_task
		;;
esac

