# opennet/packages/on-core/files/etc/hotplug.d/iface/40-on-core
#!/bin/sh

case "$ACTION" in
	ifup)
		# alle Aktionen leicht verzoegert durchfuehren, um Hotplug-Schleifen zu vermeiden
		echo "on-function update_mesh_interfaces" | on-function schedule_task
		;;
esac

