# opennet/packages/on-core/files/etc/hotplug.d/iface/40-on-core
#!/bin/sh

set -eu

case "$ACTION" in
        ifup)
		. "${IPKG_INSTROOT:-}/usr/lib/opennet/on-helper.sh"
		update_olsr_interfaces || true
		update_opennet_zone_masquerading || true
		initialize_olsrd_policy_routing || true
		;;
esac

