# opennet/packages/on-nodogsplash/files/etc/uci-defaults/70-on-captive-nodogsplash
../../usr/lib/opennet/init/on-captive-nodogsplash