# opennet/packages/on-usergw/files/etc/uci-defaults/20-on-usergw-fw-forward
../../usr/lib/opennet/init/on-usergw-fw-forward