# opennet/packages/on-usergw/files/etc/uci-defaults/on-usergw-vpn-transfer
../../usr/lib/opennet/init/on-usergw-vpn-transfer