# opennet/packages/on-usergw/files/etc/uci-defaults/on-usergw-vpn-transfer
#!/bin/sh
#
# Dieses Skript wird nur ein einziges Mal nach einem Upgrade oder der Erstinstallation ausgefuehrt:
#   http://wiki.openwrt.org/doc/uci#defaults
#

set -eu

# Allem Anschein nach tut diese Funktion folgendes:
# 1) die Konfigurationsvariable "on-usergw.ugwng_hna_mask" auf den Wert 192.168.0 setzen
# 2) openvpn-uci-Konfigurationen (z.B. openvpn.opennet_ugw_erina_on_i_de) werden nach on-usergw.opennet_ugw* übertragen
lua -e "require('luci.model.opennet.on_usergw') upgrade()" 2>/dev/null

